# mpd-client.service
[Unit]
Description=MPD Client for Waveshare 1.3inch OLED HAT
After=network.target mpd.service
Wants=mpd.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/mpd_client_fixed.py
WorkingDirectory=/usr/local/bin
Restart=always
RestartSec=10
User=root
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
